# demo/gsmwshandovercontrol.conf
[program:gsmws]
command=sudo ./GSMWSHandoverControl
directory=/home/openbts/src/gsmws
startsecs=1
stopasgroup=true

