# demo/openbts1.conf
[program:openbts1]
command=sudo ./OpenBTS
directory=/home/openbts/src/openbts-p4-ucb/rP4.0.0RC3/openbts/apps
startsecs=3
stopasgroup=true

