# demo/openbts2.conf
[program:openbts2]
command=sudo OpenBTSConfigFile=/etc/OpenBTS/OpenBTS2.db ./OpenBTS
directory=/home/openbts/src/openbts-p4-ucb/rP4.0.0RC3/openbts/apps
startsecs=15
stopasgroup=true


