# demo/sipauthserv.conf
[program:sipauthserv]
command=sudo ./sipauthserve
directory=/home/openbts/src/openbts-p4-ucb/rP4.0.0RC3/subscriberRegistry/apps

